# services.d/azuredevops-agent/run
#!/usr/bin/with-contenv sh
cd /azp
exec /azp/start.sh
